# slack-pfp-web.service
[Unit]
Description=Slack PFP Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=christian
WorkingDirectory=/home/christian/vscode/slack-pfp
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/christian/vscode/slack-pfp/.venv/bin/python /home/christian/vscode/slack-pfp/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
