# slack-pfp.service
[Unit]
Description=Slack Spotify PFP Updater
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=christian
WorkingDirectory=/home/christian/vscode/slack-pfp
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/christian/vscode/slack-pfp/.venv/bin/python /home/christian/vscode/slack-pfp/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
